# debian/ipfs.upstart
description "ipfs daemon"
author "Sergio Rubio"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 5 60

env HOME=/var/lib/ipfs

pre-start script
  if [ -f "/etc/default/ipfs" ]; then
    . /etc/default/ipfs
  fi

  if [ "$ENABLED" = "0" ]; then
    echo "WARN: ipfs daemon disabled via /etc/default/ipfs"
    stop
    exit 0
  fi
end script

exec start-stop-daemon --start -c ipfs --exec /usr/bin/ipfs -- daemon --init $DAEMON_OPTS
